# lyra_profile/airootfs/etc/profile.d/lyra-gaming.sh
export STEAM_FRAME_FORCE_CLOSE=1
export __GL_SHADER_DISK_CACHE=1
export __GL_SHADER_DISK_CACHE_SKIP_CLEANUP=1
